# cookbooks/mongodb/files/default/mongodb
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-admin/monit/files/monit.initd-4.10.1,v 1.1 2008/02/06 19:47:56 caleb Exp $

opts="${opts} reload"

depend() {
	use net
}

start() {
	ebegin "Starting mongodb"
	start-stop-daemon --start --quiet --exec /usr/local/mongodb/bin/mongod run > /dev/null 2>&1 &
	eend $?
}

stop() {
	ebegin "Stopping mongodb"
	start-stop-daemon --stop --quiet --exec /usr/local/mongodb/bin/mongod
	eend $?
}

