# cookbooks/mongodb/files/default/mongodb
#!/sbin/runscript

depend() {
  use net
}

start() {
  ebegin "Starting mongodb"
  /usr/local/mongodb/bin/mongod --master --port 27017 --dbpath /db/mongodb/master/ > /dev/null 2>&1 &
  /usr/local/mongodb/bin/mongod --slave --source localhost --port 27018 --dbpath /db/mongodb/slave/ > /dev/null 2>&1 &
  eend $?
}

stop() {
  ebegin "Stopping mongodb"
  start-stop-daemon --stop --quiet --exec /usr/local/mongodb/bin/mongod
  eend $?
}