# cookbooks/papertrail/templates/default/remote_syslog.initd.erb
#!/sbin/runscript
# Distributed under the terms of the Apache License, Version 2.0

depend() {
    need net
}

start() {
    ebegin "Starting ${NAME}"

    start-stop-daemon --start --pidfile $PIDFILE --name $NAME \
        --exec /usr/local/bin/remote_syslog -- $DAEMON_ARGS

    eend $?
}

stop() {
    ebegin "Stopping ${NAME}"

    start-stop-daemon --stop --pidfile $PIDFILE --name $NAME \
        --exec /usr/local/bin/remote_syslog

    eend $?
}
