# cookbooks/le/templates/default/logentries.initd.erb
#!/sbin/runscript

SERVICENAME="logentries"
PIDFILE="/var/run/$SERVICENAME.pid"
EXEC="/usr/bin/le-monitordaemon"

depend() {
	need net
	need localmount
	after bootmisc
}

start() {
	ebegin "Starting Logentries Agent"
	start-stop-daemon --start --quiet --oknodo --pidfile $PIDFILE --exec $EXEC
	eend $? 
}

stop() {
	ebegin "Stopping Logentries Agent"
	start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE
	eend $?
}
